# packaging/ac.socket
[Socket]
ListenStream=/run/user/%U/-2
DirectoryMode=0777

[Install]
WantedBy=sockets.target
