# packaging/ac.service
#
# classical example of a socket activated service?
#

[Unit]
Description=Start the Access Control server

[Service]
ExecStartPre=-/bin/mkdir -m 1777 /tmp/alaunch
ExecStart=/usr/bin/ac_daemon

OOMScoreAdjust=-100
IOSchedulingClass=realtime
IOSchedulingPriority=1
CPUSchedulingPriority=50

[Install]
WantedBy=multi-user.target
