# packaging/ac.socket
[Socket]
ListenStream=%t/alaunch/-2

[Install]
WantedBy=sockets.target
