# packaging/amd_session_agent.socket
[Socket]
ListenStream=/run/user/%U/amd_agent

Service=amd_session_agent.service

[Install]
WantedBy=sockets.target
