# packaging/amd_session_agent.socket
[Socket]
ListenStream=/run/user/%U/amd_agent

[Install]
WantedBy=sockets.target
