# packaging/launchpad-preload.socket
[Socket]
ListenStream=%t/alaunch/-1

[Install]
WantedBy=sockets.target
